# vms/freebsd-12.2-stable/rc.local

#!/bin/sh

# CDE
/usr/local/dt/bin/dtlogin
